# files/etc/uci-defaults/90_my_dhcp
#!/bin/sh

# Read DHCP host credentials from U-Boot env, fall back to defaults
PI5_MAC=$(fw_printenv -n dhcp_host_pi5_mac 2>/dev/null)
PI5_DUID=$(fw_printenv -n dhcp_host_pi5_duid 2>/dev/null)

if [ -z "$PI5_MAC" ]; then
    echo "WARNING: dhcp_host_pi5_mac not found in U-Boot env, using default." >&2
    PI5_MAC="00:11:22:33:44:55"
fi

if [ -z "$PI5_DUID" ]; then
    echo "WARNING: dhcp_host_pi5_duid not found in U-Boot env, using default." >&2
    PI5_DUID="000000000000000000000000000000000000"
fi

has_host() {
    local name="$1" i=0
    while uci -q get "dhcp.@host[$i]" >/dev/null 2>&1; do
        [ "$(uci -q get "dhcp.@host[$i].name")" = "$name" ] && return 0
        i=$((i + 1))
    done
    return 1
}

has_host 'BettoPi5' || {
    uci -q add dhcp host
    uci -q set dhcp.@host[-1].name='BettoPi5'
    uci -q set dhcp.@host[-1].ip='192.168.8.5'
    uci -q add_list dhcp.@host[-1].mac="$PI5_MAC"
    uci -q set dhcp.@host[-1].duid="$PI5_DUID"
    uci -q set dhcp.@host[-1].hostid='05'
}

uci commit dhcp
exit 0
