# files/etc/uci-defaults/90_my_attendedsysupgrade
#!/bin/sh

uci set attendedsysupgrade.server.github_repo='Udeste/GL-MT6000'

uci commit attendedsysupgrade

exit 0
