# files/etc/uci-defaults/90_my_ddns
#!/bin/sh

[ -e /etc/config/ddns ] || touch /etc/config/ddns

# Read DDNS credentials from U-Boot env, fall back to defaults
DDNS_LOOKUP_HOST=$(fw_printenv -n ddns_cf_lookup_host 2>/dev/null)
DDNS_DOMAIN=$(fw_printenv -n ddns_cf_domain 2>/dev/null)
DDNS_PASSWORD=$(fw_printenv -n ddns_cf_password 2>/dev/null)

if [ -z "$DDNS_LOOKUP_HOST" ]; then
    echo "WARNING: ddns_lookup_host not found in U-Boot env, using default." >&2
    DDNS_LOOKUP_HOST="router.example.com"
fi

if [ -z "$DDNS_DOMAIN" ]; then
    echo "WARNING: ddns_domain not found in U-Boot env, using default." >&2
    DDNS_DOMAIN="router.example.com@example.com"
fi

if [ -z "$DDNS_PASSWORD" ]; then
    echo "WARNING: ddns_password not found in U-Boot env, using default." >&2
    DDNS_PASSWORD="changeme123"
fi

uci -q batch << EOI
set ddns.global='ddns'
set ddns.global.ddns_dateformat='%F %R'
set ddns.global.ddns_loglines='250'
set ddns.global.upd_privateip='0'

set ddns.Cloudflare='service'
set ddns.Cloudflare.service_name='cloudflare.com-v4'
set ddns.Cloudflare.use_ipv6='0'
set ddns.Cloudflare.enabled='1'
set ddns.Cloudflare.lookup_host='$DDNS_LOOKUP_HOST'
set ddns.Cloudflare.domain='$DDNS_DOMAIN'
set ddns.Cloudflare.username='Bearer'
set ddns.Cloudflare.password='$DDNS_PASSWORD'
set ddns.Cloudflare.use_https='1'
set ddns.Cloudflare.ip_source='web'
set ddns.Cloudflare.ip_url='http://checkip.dyndns.com'
set ddns.Cloudflare.interface='wan'
set ddns.Cloudflare.use_syslog='2'
set ddns.Cloudflare.check_interval='1'
set ddns.Cloudflare.check_unit='hours'

commit ddns
EOI

exit 0
